# scripts/coco/podvm/luks-scratch/etc/systemd/system/azure-bridge-client.service
[Unit]
Description=Azure Bridge Setup for Client Container with Hostname Resolution
After=kata-agent.service network-online.target agent-protocol-forwarder.service
Wants=kata-agent.service network-online.target
Requires=agent-protocol-forwarder.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/setup-azure-bridge-client.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target